# backend/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:postgres]
command=/usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/17/main
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:uvicorn]
command=/bin/bash -c "until pg_isready -h localhost -U %(ENV_POSTGRES_USER)s 2>/dev/null; do sleep 1; done && uv run alembic upgrade head && uv run uvicorn app.main:app --host 0.0.0.0 --port 8000"
directory=/app
autostart=true
autorestart=true
priority=20
startretries=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:startup-banner]
command=/bin/bash -c "sleep 6 && echo '' && echo '  Swagger API available on: http://localhost:8010/docs' && echo ''"
autostart=true
autorestart=false
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
